# infra/cron.d/tidio-remake-backup
0 3 * * * tidio /opt/tidio-remake/infra/backup.sh
